# .config/systemd/user/gdrive-sync.timer
[Unit]
Description=Timer for the gdrive sync

[Timer]
OnCalendar=*:3/15:13

[Install]
WantedBy=timers.target
